# PythonScript_pprefEvtVolumeSizeForPDs/RAW/RAWStatisticsForRun2015Ev1inDAS.cron
*/1 * * * * lxplus0055.cern.ch /afs/cern.ch/user/h/honi/CMSSW_7_5_5_patch3/src/PythonScript_pprefEvtVolumeSizeForPDs/RAW/bash_printppPDsinDAS_RAW.sh 
